# .crontab
@hourly anacron -s -t ~/.anacron/etc/anacrontab -S ~/.anacron/spool
